# dot_config/systemd/user/hyprland.service
[Unit]
Description=Hyprland Wayland Session
After=graphical-session.target

[Service]
ExecStart=/usr/bin/start-hyprland
Restart=no
Environment=XDG_SESSION_TYPE=wayland

[Install]
WantedBy=default.target
